# supervisord.conf
[supervisord]
nodaemon=true

[program:worker]
command=python worker.py
autostart=true
autorestart=true